# detections/splunk/purple_lab/default/savedsearches.conf
## Purple Team Home Lab — Splunk Detections (Detection-as-Code)
## INSTALL: copy this whole purple_lab/ folder to $SPLUNK_HOME/etc/apps/ and restart.
##
## PREREQS (see ../../README.md):
##   - Splunk Add-on for Microsoft Windows (Splunk_TA_windows)      -> Account_Name, Service_Name,
##     Ticket_Encryption_Type, Client_Address, Process_Command_Line, New_Process_Name ...
##   - Splunk Add-on for Sysmon (Splunk_TA_microsoft_sysmon)        -> Image, DestinationIp, GrantedAccess ...
##   - Splunk Common Information Model (Splunk_SA_CIM)              -> required by the two tstats headliners
##
## Macro-free on purpose. The tstats headliners depend only on the CIM data models;
## summariesonly=false makes them work BEFORE acceleration finishes (slower), then
## accelerate transparently once summaries exist. No rex on raw Message anywhere.

# --- HEADLINERS via CIM data models (accelerate Endpoint + Network_Traffic) ---

[Purple - LOLBin Execution (CIM tstats)]
search = | tstats summariesonly=false allow_old_summaries=true count min(_time) as firstTime max(_time) as lastTime \
    from datamodel=Endpoint.Processes \
    where Processes.process_name IN ("installutil.exe","msbuild.exe","regsvr32.exe","mshta.exe","regasm.exe","regsvcs.exe") \
        AND (Processes.process="*\\Temp\\*" OR Processes.process="*\\AppData\\*" OR Processes.process="*\\ProgramData\\*") \
    by Processes.dest Processes.parent_process_name Processes.process_name Processes.process Processes.user \
  | rename "Processes.*" AS * \
  | eval firstTime=strftime(firstTime,"%F %T"), lastTime=strftime(lastTime,"%F %T")
description = T1218.004 - signed LOLBin executing a payload from a user-writable path.
dispatch.earliest_time = -24h@h
cron_schedule = */15 * * * *
enableSched = 1
alert.severity = 4

[Purple - C2 Beacon from User-Writable Path (CIM tstats)]
search = | tstats summariesonly=false allow_old_summaries=true count values(All_Traffic.dest_port) as dest_ports \
    from datamodel=Network_Traffic.All_Traffic \
    where (All_Traffic.app="*\\Temp\\*" OR All_Traffic.app="*\\AppData\\*" OR All_Traffic.app="*\\ProgramData\\*" OR All_Traffic.app="*\\Users\\Public\\*") \
    by All_Traffic.src All_Traffic.dest All_Traffic.dest_port All_Traffic.app _time span=1m \
  | rename "All_Traffic.*" AS * \
  | stats count dc(_time) as beacons earliest(_time) as firstTime latest(_time) as lastTime by src dest dest_port app \
  | where beacons > 5
description = T1071.001 - behavioral C2. Keyed on process LOCATION + egress, never on the file name (a renamed implant still trips this).
dispatch.earliest_time = -24h@h
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 5

# --- RAW sourcetype searches - TA-extracted fields (no rex) -------------------

[Purple - Kerberoasting RC4 Service Ticket]
search = source="WinEventLog:Security" EventCode=4769 Ticket_Encryption_Type="0x17" \
    Service_Name!="krbtgt" Service_Name!="*$" \
  | stats count min(_time) as firstTime max(_time) as lastTime values(Service_Name) as services \
        by Client_Address Account_Name Ticket_Encryption_Type \
  | eval firstTime=strftime(firstTime,"%F %T"), lastTime=strftime(lastTime,"%F %T")
description = T1558.003 - RC4 (0x17) TGS for a user SPN on an AES domain = Kerberoasting.
dispatch.earliest_time = -24h@h
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 4

[Purple - Password Spray (multi-user from one source)]
search = source="WinEventLog:Security" EventCode=4771 Failure_Code="0x18" \
  | bucket _time span=5m \
  | stats dc(Account_Name) as unique_users count values(Account_Name) as targets by _time Client_Address \
  | where unique_users >= 5
description = T1110.003 - one source, many accounts, one bad attempt each. Correlate on Client_Address, not per-account counts.
dispatch.earliest_time = -24h@h
cron_schedule = */5 * * * *
enableSched = 1
alert.severity = 3

[Purple - Spray Success (4771 then 4768 same source)]
search = source="WinEventLog:Security" (EventCode=4771 OR EventCode=4768) \
  | transaction Client_Address maxspan=10m \
  | search eventcount>5 EventCode=4768 \
  | table _time Client_Address Account_Name eventcount
description = Failed spray immediately followed by a TGT issuance from the same source = compromised account. P1.
dispatch.earliest_time = -24h@h
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 5

[Purple - SAM Hive Dump via reg save]
search = source="WinEventLog:Security" EventCode=4688 New_Process_Name="*\\reg.exe" \
    Process_Command_Line="*save*" \
    (Process_Command_Line="*SAM*" OR Process_Command_Line="*SYSTEM*" OR Process_Command_Line="*SECURITY*") \
  | table _time host Creator_Process_Name Process_Command_Line
description = T1003.002 - reg.exe save of credential hives. Alert on SAM+SYSTEM+SECURITY combo for low FP.
dispatch.earliest_time = -24h@h
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 4

[Purple - DCSync from non-DC]
search = source="WinEventLog:Security" EventCode=4662 \
    (Properties="*1131f6aa-9c07-11d1-f79f-00c04fc2dcd2*" OR Properties="*1131f6ad-9c07-11d1-f79f-00c04fc2dcd2*") \
    Account_Name!="*$" Account_Name!="MSOL_*" \
  | table _time host Account_Name Properties
description = T1003.006 - replication rights requested by a user account. High-fidelity, always investigate.
dispatch.earliest_time = -24h@h
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 5

[Purple - Scheduled Task Persistence (SYSTEM + user path)]
search = source="WinEventLog:Security" EventCode=4698 Account_Name!="*$" \
  | rex "<Command>(?<task_command>[^<]+)</Command>" \
  | rex "<UserId>(?<task_runas>[^<]+)</UserId>" \
  | where like(task_command,"%\\Temp\\%") OR like(task_command,"%\\AppData\\%") OR task_runas=="S-1-5-18" \
  | table _time host Account_Name Task_Name task_command task_runas
description = T1053.005 - the two rex calls parse the task XML embedded INSIDE the 4698 body (not the Message wrapper). Non-system author + SYSTEM runas or user-writable command path.
dispatch.earliest_time = -24h@h
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 4

[Purple - LSASS Access (credential theft mask)]
search = source="WinEventLog:Microsoft-Windows-Sysmon/Operational" EventCode=10 \
    TargetImage="*\\lsass.exe" \
    GrantedAccess IN ("0x1010","0x1410","0x1438","0x143a","0x1fffff") \
    SourceImage!="*\\wininit.exe" SourceImage!="*\\MsMpEng.exe" \
  | table _time host SourceImage TargetImage GrantedAccess
description = T1003.001 - handle to lsass with a read-memory access mask.
dispatch.earliest_time = -24h@h
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 4

[Purple - Kerbrute User Enumeration (4768 error 0x6)]
search = source="WinEventLog:Security" EventCode=4768 Result_Code="0x6" \
  | bucket _time span=1m \
  | stats count dc(Account_Name) as unique_users by _time Client_Address \
  | where count > 50
description = T1087.002 - burst of KDC_ERR_C_PRINCIPAL_UNKNOWN from one source = kerbrute userenum. The detection most small SOCs miss.
dispatch.earliest_time = -24h@h
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 2

# --- HARDENING VERIFICATION (re-run after Segment 10 fixes) -------------------

[Purple - Verify Account Lockout (post FGPP)]
search = source="WinEventLog:Security" EventCode=4740 \
  | table _time Target_User_Name Caller_Computer_Name
description = Confirms FGPP lockout now fires on spray attempts (absent pre-hardening).
dispatch.earliest_time = -24h@h
enableSched = 0

[Purple - RC4 Anomaly (post AES enforcement)]
search = source="WinEventLog:Security" EventCode=4769 Ticket_Encryption_Type="0x17" \
  | stats count by Service_Name Ticket_Encryption_Type
description = After AES enforcement every 0x17 is an anomaly - threshold inverts.
dispatch.earliest_time = -24h@h
enableSched = 0
